# supervisor/jupyter.conf
[program:jupyter]
command=jupyter notebook --ip=0.0.0.0
