# supervisor/jupyter.conf
[program:jupyter]
command=jupyter lab --ip=0.0.0.0 --allow-root --no-browser
