# supervisor/tensorboard.conf
[program:tensorboard]
command=tensorboard --logdir=/notebooks/tf_logs --port 6006
